# src/create/template/user-data/cloud-config/formation-units/format-ephemeral.service
[Unit]
Description=Formats the ephemeral drive

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/wipefs -f /dev/xvdb
ExecStart=/usr/sbin/mkfs.btrfs -f /dev/xvdb
