# units/enable-docker-tcp.service
[Unit]
Description=Enable the Docker Socket for the API

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl enable docker-tcp.socket
